# roles/borg/files/borgmatic-check.timer
[Unit]
Description=Run borgmatic check

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
